# telegram-bot/davidclaude-bot.service
[Unit]
Description=David OS Telegram Bot
After=network.target

[Service]
Type=simple
User=imdavid
WorkingDirectory=/home/imdavid/davidclaude/telegram-bot
ExecStart=/home/imdavid/davidclaude/telegram-bot/node_modules/.bin/tsx --env-file=/home/imdavid/davidclaude/.env --env-file=/home/imdavid/davidclaude/telegram-bot/.env src/index.ts
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
